# hpc-keepalive.service
[Unit]
Description=Preemptive SLURM job renewal

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hpc-keepalive.sh
