# hpc-vpn-renew.service
[Unit]
Description=Auto-renew EduVPN for HPC access

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hpc-vpn-renew.sh
